# infra/droplet/systemd/polisi-indexer-placeholder.service
[Unit]
Description=Polisi indexer service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=polisi
Group=polisi
WorkingDirectory=/opt/polisigpt/repo
EnvironmentFile=/opt/polisigpt/.env
ExecStartPre=/opt/polisigpt/.venv/bin/python scraper/scripts/preflight_check.py --components indexer
ExecStart=/opt/polisigpt/.venv/bin/python -m polisi_scraper.indexer.runner --mode incremental
StandardOutput=append:/opt/polisigpt/logs/indexer.log
StandardError=append:/opt/polisigpt/logs/indexer.log

[Install]
WantedBy=multi-user.target
