# infra/droplet/systemd/polisi-indexer-placeholder.service
[Unit]
Description=Polisi indexer placeholder service (Phase 2)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=polisi
Group=polisi
WorkingDirectory=/opt/polisigpt/repo
EnvironmentFile=/opt/polisigpt/.env
ExecStart=/bin/bash -lc 'echo "Indexer placeholder: Phase 2 service not implemented yet" >> /opt/polisigpt/logs/indexer.log'

[Install]
WantedBy=multi-user.target
